# vps/vps.crontab
*/15 * * * * curl -L https://canvas.place/api/board-image > /home/x4m3/web/x4m3-rocks/timelapse/img/img-$(date -u +\%Y\%m\%d-\%H\%M).png
46 02 * * * sh /home/x4m3/timelapse/vps/vps-do-it-all.sh
0 0 * * 0 rm -rfv /home/x4m3/web/x4m3-rocks/timelapse/arc/ && mkdir /home/x4m3/web/x4m3-rocks/timelapse/arc/
