# vps/vps.crontab
*/15 * * * * sh /home/x4m3/timelapse/vps/get-canvas-image.sh >/dev/null 2>&1
46 02 * * * sh /home/x4m3/timelapse/vps/vps-do-it-all.sh
