# pi/pi.crontab
00 03 * * * sh /home/pi/wkf/pi/pi-do-it-all.sh
